# cron.cron
0 4 * * * cd /var/garage/rephoto && /opt/python2713/bin/python2.7 project/manage.py refresh_albums
0 5 * * * cd /var/garage/rephoto && /opt/python2713/bin/python2.7 project/manage.py clearsessions
15 5 * * * cd /var/garage/rephoto && /opt/python2713/bin/python2.7 project/manage.py reverse_geocode_photos
30 5 * * * cd /var/garage/rephoto && /opt/python2713/bin/python2.7 project/manage.py update_index
20 6 * * * sudo sh /var/garage/rephoto/inode-alert
