# fotodugnad.conf
[program:fotodugnad]
directory=/home/fotodugnad/ajapaik-web
command=/home/fotodugnad/venv/bin/uwsgi
  --socket=/home/fotodugnad/ajapaik-web/fotodugnad.sock
  --module=project.ajapaik.wsgihandler
  --chmod-socket=777
  --env LANG="en_US.utf8"
  --processes=8
  --master
  --harakiri=180
  --max-requests=5000
autostart=true
autorestart=true
user=fotodugnad
stopasgroup=true
environment=DJANGO_SETTINGS_MODULE="project.ajapaik.settings.local"