# solr.conf
[program:solr]
directory=/home/ajapaik/solr-4.10.4/example
command=java -jar start.jar -Djetty.host=127.0.0.1
user=ajapaik